# system-configs/sneezybackup.timer
[Unit]
Description=Back up SneezyMUD

[Timer]
OnCalendar=10:00

[Install]
WantedBy=timers.target
